# root/etc/services.d/rar2fs/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

exec s6-setuidgid hotio rar2fs -f --date-rar --seek-length=1 -o allow_other -o ro "${INPUT_DIR}" "${OUTPUT_DIR}"
